# roles/registry_maintenance/files/registry-gc.timer
[Unit]
Description=Run Docker Registry Garbage Collection Weekly

[Timer]
OnCalendar=Sun 02:00
Persistent=true

[Install]
WantedBy=timers.target
